# devbox/idle/devbox-idle.service
[Unit]
Description=Idle shutdown check for dev box

[Service]
Type=oneshot
ExecStart=/usr/local/bin/idle-check
